# rootfs/etc/udev/rules.d/11-usb-test.rules
ACTION!='add', GOTO='farsight'
KERNEL=='sd[a-z][0-9]', RUN+='/opt/usb/flashusb.sh %k' 
KERNEL=='sd[a-z]', RUN+='/opt/usb/flashusb.sh %k'
LABEL='farsight'
